# src/core/vmm/templates/glusterfs-sync.service
[Unit]
Description=Mount GlusterFS and bind mount directories
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/glusterfs-sync.sh
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/glusterfs-sync.env

[Install]
WantedBy=multi-user.target
